# nvidia.service
[Unit]
Description=NVIDIA driver
Before=dcos-mesos-slave.service

[Service]
ExecStart=/opt/nvidia/startup.sh

[Install]
WantedBy=multi-user.target
